# scripts/com.vypearr.job-agent-handshake.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.vypearr.job-agent-handshake</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/tarunva/Job_Automation/scripts/run_unattended_handshake_cycle.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/tarunva/Job_Automation</string>
  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key>
      <integer>10</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>17</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
  </array>
  <key>ProcessType</key>
  <string>Interactive</string>
  <key>StandardOutPath</key>
  <string>/Users/tarunva/Job_Automation/data/logs/launchd.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/tarunva/Job_Automation/data/logs/launchd.stderr.log</string>
</dict>
</plist>
